# service/nginx/run
#!/bin/bash -e

/usr/sbin/nginx -g 'error_log /dev/null crit;' -c /home/melpa/melpa/cfg/nginx
