# cfg/etc/runit/melpa/run
#!/bin/sh -e
exec 2>&1 
exec chpst -u melpa runsvdir /home/melpa/service 'log: ...........................................................................................................................................................................................................................................................................................................................................................................................................'
