# cfg/home/service/nginx-stable/run
#!/bin/bash -e

exec /usr/sbin/nginx -g 'error_log /dev/null crit;' -c ~/etc/nginx-stable

