# service/caddy/run
#!/bin/bash -e

/usr/local/bin/caddy --conf /home/melpa/Caddyfile
