# service/caddy/run
#!/bin/bash -e

exec /usr/local/bin/caddy --agree --log /home/melpa/log/caddy.log --conf /home/melpa/Caddyfile
